# voice-builder-api/supervisord/cherrypy.conf
[supervisord]

[program:maryttsapi]
command=gunicorn wsgi --config /opt/marytts-api/gunicorn.conf --chdir=/opt/marytts-api
