# voice-builder-api/supervisord/cherrypy.conf
[supervisord]

[program:maryttsapi]
command=gunicorn wsgi --config /opt/voice-builder-api/gunicorn.conf --chdir=/opt/voice-builder-api
