# voicebuilder/server/cherrypy/supervisord/cherrypy.conf
[supervisord]

[program:voicebuilderapi]
command=gunicorn wsgi --config /opt/marytts/voice-builder-api/gunicorn.conf --chdir=/opt/marytts/voice-builder-api
